# sv/statsd_exporter/run
#!/bin/bash

#sv start socklog-unix || exit 1

source /etc/envvars

cd /statsd_exporter

exec 2>&1
exec ./statsd_exporter
