# sv/blackbox-exporter/run
#!/bin/bash

#sv start socklog-unix || exit 1

source /etc/envvars

cd /blackbox_exporter
exec 2>&1
exec ./blackbox_exporter -config.file=blackbox.yml 
