# sv/grafana/run
#!/bin/bash

sv start prometheus || exit 1

source /etc/envvars
env

cd /grafana/bin

exec 2>&1
exec ./grafana-server
